# deploy/systemd/miodek-exporter.container
# Quadlet Podman dla eksportera metryk (KAN-219) — wariant PREFEROWANY dla homelaba na Podmanie.
# Zero-dep, więc stockowy obraz python:3-slim z zamontowanym kodem (bez budowania obrazu),
# spójnie z resztą stacku (wszystko jako quadlet, nie gołe usługi).
#
# Instalacja (rootless, user): ~/.config/containers/systemd/miodek-exporter.container
#   1. rsync kodu na host do ~/miodek (bez .git/tests), korpus do ~/miodek-corpus
#   2. systemctl --user daemon-reload && systemctl --user start miodek-exporter.service
#   3. curl localhost:9112/metrics
# Wariant gołej usługi systemd: patrz miodek-exporter.service (gdy host bez Podmana).

[Unit]
Description=Sztuczny Miodek metrics exporter (KAN-219)
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/library/python:3-slim
AutoUpdate=registry
Network=host
Volume=%h/miodek:/app:ro,Z
Volume=%h/miodek-corpus:/corpus:ro,Z
Environment=MIODEK_CORPUS=/corpus
Environment=MIODEK_PORT=9112
Environment=MIODEK_LOG=/app/decisions.jsonl
# Eksporter żyje w pakiecie miodek (src/). Shim tools/metrics_exporter.py sam dokłada /app/src,
# ale ustawiamy też PYTHONPATH jako zabezpieczenie — wymaga pełnego klonu repo z katalogiem src/.
Environment=PYTHONPATH=/app/src
Exec=python3 /app/tools/metrics_exporter.py --port 9112 --corpus /corpus

[Service]
Restart=always

[Install]
WantedBy=default.target
